# deploy/com.spectral160.autodeploy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.spectral160.autodeploy</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/serverpibes/macmini-vps/docker/spectral160-presskit/deploy/auto-deploy.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>120</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/spectral160-autodeploy.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/spectral160-autodeploy.err</string>
</dict>
</plist>
